# redmine-n8n/mem-agent-mcp/com.musai.memagent.wrapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.musai.memagent.wrapper</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/christopherhicks/redmine-n8n/mem-agent-mcp/bin/serve_mcp_wait.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/christopherhicks/.local/bin</string>
    <key>LMSTUDIO_BASE_URL</key><string>http://127.0.0.1:8000</string>
    <key>MEMORY_DIR</key><string>/Users/christopherhicks/redmine-n8n/data/notes/MusaiMemory</string>
    <key>MEMORY_PATH</key><string>/Users/christopherhicks/redmine-n8n/data/notes/MusaiMemory</string>
    <key>MEM_AGENT_MEMORY_DIR</key><string>/Users/christopherhicks/redmine-n8n/data/notes/MusaiMemory</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/christopherhicks/redmine-n8n/mem-agent-mcp</string>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>StandardOutPath</key><string>/Users/christopherhicks/Library/Logs/memagent-mcp.out</string>
  <key>StandardErrorPath</key><string>/Users/christopherhicks/Library/Logs/memagent-mcp.err</string>
</dict>
</plist>
